# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:node-socket]
command=node /usr/local/app/socket/index.js

[program:python-scraper]
command=python3 /usr/local/app/app.py
